# etc/profile.d/banner.sh
#!/bin/sh

echo "--------------------------------"
cat << EOF
"
в Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв ҖвЈҖвЈҖвЎҖв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җ
в Җв Җв Җв Җв Җв Җв Җв ҖвўҖв Җв Җв ҖвЎҮв Җв Ғв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җ
в Җв Җв Җв Җв Җв Җв Җв ҖвўҢв ҷв ўвЈҖвЎҮв ҖвЈҖв Өв Ҡвўҷв „в Җв Җв Җв Җв Җв Җв Җ
в Җв Җв Җв Җв Җв үв үвЎҖв Җв ‘в ўвЈҖв Ғв ҖвўҒвЎ в Ҡв Ғв ҖвЎҢвў©в Җв Җв Җв Җв Җ
вўҖвЎһв ‘в ўвў„вў°в ҖвЎҮв Җв Җв Җв ҲвЎҮв Җв Ғв Җв Җв Җв ҖвЎҮвўёвЈҖв Өв Ҡв ҷвЈ„
в Җв үв Ӯвў„вЎҖв үв ҖвўҮвЎҖв Җв Җв ҖвЎҮв Җв Җв Җв Җв ҖвЈ в Үв Ҳв ҒвЎ в ”в Ҡв Җ
в ҖвЎ в ”в ӮвўҒвЎ в ӨвЈҖв Ҳв ’вў„вЎҖвЎҮв Җв ҖвЎҖв ”в ҠвўҖв  в ӨвЈҲв ‘в Ӯв ӨвЎҖ
в Җв ұв ”в Ҡв Ғв Җв Җв Җв үв ўвў„вЎҲв Үв Җв үвЈҖв „в Ҡв Ғв Җв Җв Җв үв ’в ”в Ғ
в ҖвЎ в ӨвЈҖв Җв Җв Җв ҖвўҖвЎ в җв үвЎҖв ҖвЎҲв ‘в ӨвЈҖв Җв Җв Җв ҖвўҖвЎ вўӨв Җ
в Җв “в ӨвЈҖв Ғв Ӯв Өв Ӯв ҒвЎ в ”в үвЎҮв Җв үв ’в ӨвЎҖв Ғв ўв ”в ҠвўҒвЎ в җв ғ
в ҖвўҖвЎ в ”в үвўҖв ҖвЎ–в үв Җв Җв ҖвЎҮв Җв Җв Җв Җв Ҳв ‘вЎ„вўҖв Ҳв “в ӨвЎҖв Җ
в җвўҮвўҖв Өв ’в №в ҖвЎҮв Җв Җв Җв ҖвЎҮв ҖвЎҖв Җв Җв Җв ҖвЎҮвўёв ‘в ўвЈҖвўҲв Ү
в Җв Ҳв Ғв Җв ҖвЎ„вўҖв Үв ҖвўҖв Өв ҠвЎҒв Җв үв ўвў„в Җв ҖвЈҮвўҳв Җв Җв Җв үв Җ
в Җв Җв Җв Җв Җв үв Ғв ҖвўҺвўҒв Өв ҠвЎҮв Җв үв ўвў„вўұв Ӯв Җв үв Җв Җв Җв Җв Җ
в Җв Җв Җв Җв Җв Җв Җв Җв Ҳв Ғв Җв ҖвЎҮв ҖвЎҖв Җв Җв үв Җв Җв Җв Җв Җв Җв Җв Җ
в Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Ғв Җв Ғв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җв Җ
OS
"
EOF
echo "--------------------------------"
echo "SnowOS v1.0.0 (Snowfall) based on Ubuntu 26.04 LTS (Resolute Racoon)"
echo
echo "GitHub: SnowOS-Linux"
echo "Report bugs: https://github.com/SnowOS-Linux/SnowOS/issues"
echo "This is not a Canonical official flavor. If you want official Canonical releases or flavors, see: https://ubuntu.com/desktop/flavors"
echo "SnowOS It does not set a default password. Use the passwd command to set a password."
echo 
echo "Using the Package Manager"
echo "Install a package: spkg install <pkg>"
echo "Update repositories: spkg update"
echo "Upgrade: spkg upgrade"
echo "Download a .deb package: spkg download <pkg>"
echo "Current user: $(whoami)"
echo
